# deploy/gunicorn/gunicorn_http.service
[Unit]
Description=Gunicorn HTTP daemon for LitLoop
After=network.target

[Service]
User=driptamine
Group=www-data
WorkingDirectory=/home/driptamine/litloop_backend_v2
Environment=DJANGO_SETTINGS_MODULE=litloop_project.settings.prod

# Run with multiple workers for better HTTP performance
ExecStart=/home/driptamine/litloop_backend_v2/env/bin/gunicorn \
    --preload \
    --workers 2 \
    --max-requests 10000 \
    --max-requests-jitter 1000 \
    --timeout 60 \
    --bind 127.0.0.1:8001 \
    litloop_project.wsgi:application

Restart=always

[Install]
WantedBy=multi-user.target
