# deploy/gunicorn/gunicorn.service
[Unit]
Description=Daphne ASGI daemon for LitLoop (WebSockets only)
After=network.target

[Service]
User=driptamine
Group=www-data
WorkingDirectory=/home/driptamine/litloop_backend_v2
Environment=DJANGO_SETTINGS_MODULE=litloop_project.settings.prod
LimitNOFILE=65536

# Dedicated Daphne for WebSockets on port 8000
ExecStart=/home/driptamine/litloop_backend_v2/env/bin/daphne \
    -v 2 \
    -b 127.0.0.1 \
    -p 8000 \
    litloop_project.asgi:application

Restart=always

[Install]
WantedBy=multi-user.target

